# project_name/conf/livesite.conf
description "Django site"
start on runlevel [2345]
stop on runlevel [06]
respawn
respawn limit 10 5
exec /var/sites/virtualenvs/{{ project_name }}/src/{{ project_name }}/conf/launch.sh
